# dashboard/trading-dashboard.service
[Unit]
Description=Trading System Status Dashboard
After=network.target

[Service]
Type=simple
User=mirko
WorkingDirectory=/opt/trading/traderunner/dashboard
ExecStart=/usr/bin/python3 /opt/trading/traderunner/dashboard/server.py
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
